# deploy/systemd/pain-farm.service
[Unit]
Description=Pain Farm dashboard (Next.js production)
After=network-online.target core-swicht.service
Wants=core-swicht.service

[Service]
Type=simple
User=server
Group=server
WorkingDirectory=/home/server/pain-farm
Environment=NODE_ENV=production
Environment=BACKEND_URL=http://127.0.0.1:8000
Environment=NEXT_PUBLIC_BOOT_RUN_MODE=main
Environment=HOSTNAME=0.0.0.0
Environment=PORT=3000
ExecStart=/bin/bash -lc 'export NVM_DIR="/home/server/.nvm"; . "/home/server/.nvm/nvm.sh"; cd /home/server/pain-farm && npm run start'
Restart=always
RestartSec=8

[Install]
WantedBy=multi-user.target
